# http/kvm/kvm-hostname.service
[Unit]
Description=KVM hostname

[Service]
Type=oneshot
ExecStart=/usr/local/bin/kvm-hostname.sh

[Install]
WantedBy=multi-user.target
